# rvm/soft/supervisord.d/jekyll.conf
[program:jekyll]
command=/root/start.sh
priority=700
